# tunnel-service/autossh.service
[Unit]
Description= Automatic SSH tunnel

[Service]
ExecStart=/usr/bin/ssh :q

